# proxyServer.service
[Unit]
Description=Proxy Server Service
After=apache2.service
Requires=apache2.service

[Service]
Type=simple
ExecStart=/usr/bin/socat TCP4-LISTEN:80,fork,pktinfo system:/etc/proxyServer.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
